# rpm/before-remove.sh
#!/bin/bash
/etc/init.d/beamium stop
chkconfig --del beamium
