# rpm/before-remove.sh
#!/bin/bash
/etc/init.d/noderig stop
chkconfig --del noderig
