# deb/before-remove.sh
#!/bin/bash
/etc/init.d/noderig stop
update-rc.d -f noderig remove
